# systemd.service
[Unit]
Description=InstagramBot (aiogram) service
After=network.target

[Service]
Type=simple
User=botuser
Group=botuser
WorkingDirectory=/srv/instagram-bot
Environment="PATH=/srv/instagram-bot/venv/bin"
# If you use a virtualenv, point to its python
ExecStart=/srv/instagram-bot/venv/bin/python /srv/instagram-bot/main.py
Restart=on-failure
RestartSec=5
# limit resources (tune as needed)
LimitNOFILE=65536
# Optional: log to journald (systemd handles stdout/stderr)
StandardOutput=inherit
StandardError=inherit

[Install]
WantedBy=multi-user.target

#sudo systemctl daemon-reload
#sudo systemctl enable instagram-bot.service
#sudo systemctl start instagram-bot.service
#sudo journalctl -u instagram-bot.service -f
